# contrib/systemd/peca-p2p-yp.service
[Unit]
Description=peca-p2p-yp - decentralized P2P yellow pages node
After=network-online.target
Wants=network-online.target

[Service]
Type=notify
ExecStart=/usr/local/bin/peca-p2p-yp
User=peca-p2p-yp
Group=peca-p2p-yp
StateDirectory=peca-p2p-yp
StateDirectoryMode=0700
UMask=0077
Restart=on-failure
RestartSec=5

# ハードニング(最小権限 — Principle II)
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
RestrictNamespaces=true
LockPersonality=true
MemoryDenyWriteExecute=true
RestrictRealtime=true
RestrictSUIDSGID=true
SystemCallArchitectures=native
SystemCallFilter=@system-service
CapabilityBoundingSet=

[Install]
WantedBy=multi-user.target
